# data/systemd/clad@.socket
[Unit]
Description=Command Line Assistant Daemon Socket

[Socket]
ListenStream=clad@_%i

[Install]
WantedBy=sockets.target
